# config/systemd/user/vlx-bl-tracker.service
[Unit]
Description=Bundesliga live score check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vlx bl setup --poll
